# abovetopsecret-dash/deploy/opticdata-webhook.service
[Unit]
Description=OpticData GitHub Webhook Listener
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/abovetopsecret-dash/abovetopsecret-dash/deploy
ExecStart=/usr/bin/node /opt/abovetopsecret-dash/abovetopsecret-dash/deploy/webhook-listener.js
Restart=always
RestartSec=5
EnvironmentFile=/opt/abovetopsecret-dash/abovetopsecret-dash/deploy/.env.webhook

[Install]
WantedBy=multi-user.target
